# templates/vps/cron/backup.cron.example
# Daily backup example at 02:00
0 2 * * * /home/gitnanu/scripts/backup.sh
